# etc/systemd/fpp-announce-ip.service
[Unit]
Description=FPP Announce IP Address
# Run after the network setup has resolved an IP and after PipeWire is up so the
# flite-rendered announcement can play. Deliberately NOT ordered before fppd:
# this is fire-and-forget so the (CPU-heavy on single-core boards) flite render
# no longer competes with the tail of fpp_postnetwork or fppd startup.
After=fpp_postnetwork.service fpp-pipewire.service fpp-wireplumber.service fpp-pipewire-pulse.service
Wants=fpp-pipewire.service fpp-wireplumber.service

[Service]
Type=oneshot
ExecStart=/opt/fpp/src/fppinit announceIP

[Install]
WantedBy=multi-user.target
