# etc/systemd/fpp-early-block-trigger.service
[Unit]
Description=FPP - Early coldplug of block and sound devices
Documentation=https://github.com/FalconChristmas/fpp
# The kernel creates the mmc/SD partition nodes (and, on capes whose audio codec
# is in the boot device tree, the sound card) very early (~4s), but on a
# single-core SBC systemd's full udev coldplug processes ~700 AM335x devices
# serially, so those lightweight uevents aren't tagged until ~28s. For block
# that means every local mount (/, /boot/firmware, swap) -- and the FPP services
# ordered after local-fs.target/swap.target -- waits it out. For sound it means
# WirePlumber sees the card late, stalling fpp_postnetwork's audio-device wait.
# Trigger and settle just these two subsystems first, before the full coldplug,
# so they resolve in ~1s while the rest coldplugs in the background. (On a
# no-soundcard board snd-dummy isn't loaded yet, so the sound match is a no-op.)
DefaultDependencies=no
Wants=systemd-udevd.service
After=systemd-udevd-kernel.socket systemd-udevd-control.socket systemd-udevd.service
Before=systemd-udev-trigger.service local-fs-pre.target shutdown.target
Conflicts=shutdown.target

[Service]
Type=oneshot
RemainAfterExit=yes
# Block first, WITH --settle: fppoled/fppinit and the local mounts need these
# resolved right now, so wait for them. The settle is fast because only block
# is queued at this point. Then trigger sound WITHOUT --settle: we don't need
# the card resolved until fpp_postnetwork (much later), so don't block the
# critical path on it -- just enqueue it ahead of the full coldplug so it
# processes early and is ready by then. (Order matters: settling block before
# enqueuing sound keeps the settle from waiting on sound.)
ExecStart=/usr/bin/udevadm trigger --action=add --subsystem-match=block --settle
ExecStart=/usr/bin/udevadm trigger --action=add --subsystem-match=sound

[Install]
WantedBy=sysinit.target
