# etc/systemd/fppd.service
[Unit]
Description=FPPd
After=network.target generic-board-startup.service remote-fs.target getty.target fpp_postnetwork.service
Requires=network.target 
DefaultDependencies=no

[Service]
Type=simple
ExecStart=/opt/fpp/src/fppd -f
ExecStop=/opt/fpp/scripts/fppd_stop
ExecStartPre=-/opt/fpp/scripts/fppd_boot_pre
ExecStartPost=-/opt/fpp/scripts/fppd_boot_post
ExecReload=/opt/fpp/src/fpp -r
Restart=always


[Install]
WantedBy=multi-user.target
