# etc/systemd/fppd.service
[Unit]
Description=FPPd
After=network.target apache2.target generic-board-startup.service
Requires=network.target 
Want=generic-board-startup.service
DefaultDependencies=no


[Service]
Type=oneshot
RemainAfterExit=yes
ExecStart=/opt/fpp/scripts/fppd_start start
ExecStop=/opt/fpp/scripts/fppd_stop


[Install]
WantedBy=network.target multi-user.target
