# etc/systemd/fppinit.service
[Unit]
Description=Inialization routines for FPP
DefaultDependencies=no
Before=network.service fpprtc.service fppoled.service home-fpp-media.mount
After=home-fpp-media.mount

[Service]
Type=oneshot
RemainAfterExit=yes
ExecStartPre=-mkdir -p /home/fpp/media/logs
ExecStart=/opt/fpp/src/fppinit start
StandardOutput=append:/home/fpp/media/logs/fpp_init.log
StandardError=append:/home/fpp/media/logs/fpp_init.log

[Install]
WantedBy=sysinit.target
